# docker/s6/s6-rc.d/nextjs/run
#!/bin/sh
export PORT=3000
export HOSTNAME=127.0.0.1
cd /app/web
exec node server.js
